# dnssec-tools/dist/donutsd.service.in
[Unit]
Description=Donuts DNS lint checking daemon
After=network.target

[Service]
ExecStart=@prefix@/bin/donutsd
Type=simple

[Install]
WantedBy=multi-user.target
