# dnssec-tools/dist/rollerd.service.in
[Unit]
Description=Rollerd DNSSEC rolling daemon
After=network.target

[Service]
ExecStart=@prefix@/bin/rollerd -foreground
Type=simple

[Install]
WantedBy=multi-user.target
